# debian/postrm
#!/bin/sh
set -e
if [ "$1" = "remove" ]; then
    ldconfig
    update-rc.d -f dabba remove
fi
